# dot_config/systemd/user/backlight.service
[Unit]
Description=Run vimbacklight script

[Service]
ExecStart=%h/.config/systemd/user/vimbacklight.sh

[Install]
WantedBy=default.target
