# dot_config/systemd/user/backlight.timer
[Unit]
Description=Run vimbacklight script 1 minute after login

[Timer]
OnActiveSec=1min
Unit=backlight.service

[Install]
WantedBy=default.target
