# dot_config/systemd/user/timers.target.wants/epub.timer
[Unit]
Description=Run epub service once a week

[Timer]
OnBootSec=15min
OnUnitActiveSec=5h
Unit=epub.service

[Install]
WantedBy=timers.target
